# conf/supervisor/scrollcontrol.conf
[program:scrollcontrol]
command=/usr/bin/node /home/pi/scrollcontrol/server/server.js
autorestart=true
user=pi
directory=/home/pi/scrollcontrol/server
